# scripts/phone-server/units/daanio-serve.service
[Unit]
Description=daanio server with WebSocket gateway
After=network.target

[Service]
Environment=PATH=/home/ec2-user/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=DAANIO_BEDROCK_ENABLE=1
Environment=AWS_REGION=us-east-1
Environment=DAANIO_BEDROCK_MODEL=us.anthropic.claude-opus-4-6-v1
ExecStartPre=/bin/rm -f /run/user/1000/daanio-daemon.lock /run/user/1000/daanio.sock /run/user/1000/daanio.sock.hash /run/user/1000/daanio-debug.sock
ExecStart=/home/ec2-user/.local/bin/daanio --provider bedrock serve
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
